# demo_docker/atramhasis.conf
[supervisord]
nodaemon=true

[program:atramhasis]
command=pserve /opt/atramhasis_demo/development.ini

[program:ldf_server]
command=ldf-server /opt/atramhasis_demo/ldf_server_config.json
